# supervisord.cloudrun.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
user=root

[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

# ── uvicorn (FastAPI backend on internal port 8001) ────────────────────────
[program:uvicorn]
command=python -m uvicorn main:app --host 127.0.0.1 --port 8001 --workers 1 --timeout-keep-alive 900
directory=/app
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=HOME="/root",HF_HUB_DISABLE_SYMLINKS_WARNING="1"

# ── nginx (reverse proxy + static file server on $PORT) ───────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
# nginx must start after uvicorn is ready to accept connections.
startsecs=3
